# Additive/vgw/etc/systemd/system/additive.service.j2
[Unit]
Description = {{ description }}
After = network.target

[Service]
Type = oneshot
User = root
ExecStart = /bin/sh /usr/bin/additive.sh 

[Install]
WantedBy = multi-user.target